# crontabs
0 */1 * * 1-5 /home/kiosk/cron-kiosk.sh
0 */1 * * 1-5 /home/kiosk/cron-kiosk-update.sh
